# firewall/service/firewall.service
[Unit]
Description=Custom Firewall Manager
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/firewall/manager/main.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
